# com.quentin.dictate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quentin.dictate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/quentinlachaud/code/qispr/.venv/bin/python</string>
        <string>/Users/quentinlachaud/code/qispr/dictate.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DICTATE_INPUT_DEVICE</key>
        <string>Logitech Webcam</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>LimitLoadToSessionType</key>
    <array>
        <string>Aqua</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/quentinlachaud/Library/Logs/dictate.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/quentinlachaud/Library/Logs/dictate.stderr.log</string>
</dict>
</plist>
